# crontab
# Algo Trading Crontab
# Install with: crontab /home/jay/dev/algo/crontab
# All times are ET (Eastern Time)

# Backfill outcomes (6 AM, before market open, Mon-Fri)
0 8 * * 1-5 /home/jay/dev/algo/run-docker.sh trading python -m trading.backfill >> /home/jay/dev/algo/logs/backfill.log 2>&1

# News pipeline (6:30 AM, before ideation, Mon-Fri)
30 8 * * 1-5 /home/jay/dev/algo/run-docker.sh trading python -m trading.pipeline >> /home/jay/dev/algo/logs/pipeline.log 2>&1

# Ideation - generate/review theses (7 AM, after pipeline, Mon-Fri)
0 9 * * 1-5 /home/jay/dev/algo/run-docker.sh trading python -m trading.ideation_claude >> /home/jay/dev/algo/logs/ideation.log 2>&1

# News pipeline (4 PM, before trading session, Mon-Fri)
0 16 * * 1-5 /home/jay/dev/algo/run-docker.sh trading python -m trading.pipeline >> /home/jay/dev/algo/logs/pipeline.log 2>&1

# Trading session (4:30 PM, 30 min after market close, Mon-Fri)
30 18 * * 1-5 /home/jay/dev/algo/run-docker.sh trading python -m trading.trader >> /home/jay/dev/algo/logs/trader.log 2>&1

# Strategy evolution (7 AM Sunday, weekly)
0 7 * * 0 /home/jay/dev/algo/run-docker.sh trading python strategy.py --days 60 >> /home/jay/dev/algo/logs/strategy.log 2>&1
