# crontab
# Algo Trading Crontab
# Install with: crontab /home/jay/dev/algo/crontab
# All times are ET (Eastern Time)

# Backfill outcomes (6 AM ET, Mon-Fri) — learning system
0 8 * * 1-5 /home/jay/dev/algo/run-docker.sh trading python -m trading.backfill >> /home/jay/dev/algo/logs/backfill.log 2>&1

# Consolidated daily session (3 PM ET, Mon-Fri)
0 19 * * 1-5 /home/jay/dev/algo/run-docker.sh trading python -m trading.session >> /home/jay/dev/algo/logs/session.log 2>&1

# Learning loop (7 AM ET Sunday) — weekly learning
0 7 * * 0 /home/jay/dev/algo/run-docker.sh trading python -m trading.learn --days 60 >> /home/jay/dev/algo/logs/learn.log 2>&1
